# backend/systemd/samsung-auto-update.timer
[Unit]
Description=Run Samsung Display Hub Auto Update every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=20
Persistent=true

[Install]
WantedBy=timers.target
